# cron/crontab
*/30	*	*	*	*	cat /dev/null > ~/.zsh_history
*/30	*	*	*	*	cp ~/Library/Safari/Bookmarks.plist ~/Dropbox/Library/Bookmarks
